# scripts/launchd/com.konecta.contadores.crm-followup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.konecta.contadores.crm-followup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/fgoiriz/private/repos/contadores/scripts/run_contadores_crm_hourly_followup.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/fgoiriz/private/repos/contadores</string>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/fgoiriz/private/repos/contadores/data/reports/launchd-contadores-crm-followup.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/fgoiriz/private/repos/contadores/data/reports/launchd-contadores-crm-followup.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/fgoiriz</string>
    <key>CODEX_HOME</key>
    <string>/Users/fgoiriz/.codex</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
